# .config/autostart/mpd.desktop
[Desktop Entry]
Name=mpd
Exec=mpd --no-daemon
Terminal=false
Type=Application
