# .config/autostart/xbanish.desktop
[Desktop Entry]
Name=xbanish
Exec=xbanish
Terminal=false
Type=Application
