# .config/autostart/xss-lock.desktop
[Desktop Entry]
Name=xss-lock
Exec=xss-lock xlock
Terminal=false
Type=Application
