# sv/kdeconnectd/run
#!/bin/sh
exec /usr/lib/libexec/kdeconnectd

